# files/etc/crontabs/root
6 4 * * 3 /usr/bin/apk update && /usr/bin/apk upgrade
